# CentOS/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:mariadb]
command=/usr/libexec/mysqld --datadir=/var/lib/mysql --user=mysql
autostart=true
autorestart=true
stderr_logfile=/var/log/mariadb/mariadb_error.log
stdout_logfile=/var/log/mariadb/mariadb.log

[program:sshd]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/sshd_err.log
stdout_logfile=/var/log/supervisor/sshd_out.log

[program:php-fpm]
command=/usr/sbin/php-fpm --nodaemonize
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/php-fpm_err.log
stdout_logfile=/var/log/supervisor/php-fpm_out.log

[program:httpd]
command=/usr/sbin/httpd -DFOREGROUND
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/httpd.err.log
stdout_logfile=/var/log/supervisor/httpd.out.log