# MySQL/8.x/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:mysql]
command=/usr/sbin/mysqld
autostart=true
autorestart=true
